# launchd/com.brand-deal-agent.schedule.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.brand-deal-agent.schedule</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/khristina/.nvm/versions/node/v22.6.0/bin/node</string>
    <string>src/schedule.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/khristina/Documents/Brand Deal Agent</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/khristina/Library/Logs/brand-deal-agent.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/khristina/Library/Logs/brand-deal-agent-error.log</string>
</dict>
</plist>
